# containers/quadlet/oxidized.container
# Podman Quadlet configuration for Oxidized
# Install this file to: /etc/containers/systemd/oxidized.container
# Then run: systemctl daemon-reload && systemctl enable --now oxidized
#
# This file generates a systemd service that manages the Oxidized container
# with automatic restarts and persistent storage.
#
# Security model: Runs as dedicated oxidized user (UID 2000, GID 2000)
# with enhanced security constraints (read-only rootfs, dropped capabilities)

[Unit]
Description=Oxidized Network Configuration Backup Service
Documentation=https://github.com/yggdrasil-network/oxidized
After=network-online.target
Wants=network-online.target

[Container]
# Pinned to stable release (not latest)
# Check https://hub.docker.com/r/oxidized/oxidized/tags for versions
Image=docker.io/oxidized/oxidized:0.30.1

# Container name for easier management
ContainerName=oxidized

# Network: Dedicated bridge network (no host networking needed)
Network=oxidized-net

# Security: Run as dedicated non-root user (oxidized:oxidized = 2000:2000)
User=2000:2000

# Security: Read-only root filesystem
ReadOnly=true

# Security: Drop all capabilities (Oxidized doesn't need any)
DropCapability=ALL

# Security: Prevent privilege escalation
SecurityLabelDisable=false
NoNewPrivileges=true

# Writable tmpfs for runtime needs
Tmpfs=/tmp:rw,noexec,nosuid,size=64m
Tmpfs=/run:rw,noexec,nosuid,size=32m

# Persistent storage with SELinux relabeling (:Z for exclusive container access)
# Host -> Container mappings (all owned by 2000:2000 on host)
Volume=/var/lib/oxidized/config:/home/oxidized/.config/oxidized:Z
Volume=/var/lib/oxidized/ssh:/home/oxidized/.ssh:Z,ro
Volume=/var/lib/oxidized/data:/home/oxidized/.config/oxidized/data:Z
Volume=/var/lib/oxidized/output:/home/oxidized/.config/oxidized/output:Z
Volume=/var/lib/oxidized/repo:/home/oxidized/.config/oxidized/repo:Z

# Environment variables
Environment=TZ=UTC
Environment=HOME=/home/oxidized

# Health check via netcat (curl not available in read-only rootfs without /tmp)
# Using simple TCP check instead
# Note: Health checks may not work with read-only rootfs; monitor via systemd

[Service]
# Restart policy for production reliability
Restart=always
RestartSec=10
TimeoutStartSec=120
TimeoutStopSec=30

# Resource limits (adjust based on device count)
# For ~100 devices, these are conservative
MemoryLimit=1G
CPUQuota=100%

[Install]
# Enable service to start on boot
WantedBy=multi-user.target default.target
